# templates/jaeger-query.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=jaeger-query
After=network-online.target
StartLimitInterval=0
StartLimitIntervalSec=0

[Service]
Type=simple
PIDFile=/var/run/jaeger-query.pid
User={{ jaeger_system_user }}
Group={{ jaeger_system_group }}
PermissionsStartOnly=true
Environment="SPAN_STORAGE_TYPE={{ jaeger_query_span_storage_type }}"
ExecReload=/bin/kill -HUP $MAINPID
ExecStart=/usr/local/bin/jaeger-query --config-file /etc/jaeger-query.yml

SyslogIdentifier=jaeger-query
KillMode=process
Restart=always
RestartSec=5

LockPersonality=true
NoNewPrivileges=true
MemoryDenyWriteExecute=true
PrivateTmp=true
ProtectHome=true
RemoveIPC=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
