# templates/jaeger-agent.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=jaeger-agent
After=network-online.target
StartLimitInterval=0
StartLimitIntervalSec=0

[Service]
Type=simple
PIDFile=/var/run/jaeger-agent.pid
User={{ jaeger_system_user }}
Group={{ jaeger_system_group }}
PermissionsStartOnly=true
Environment="SPAN_STORAGE_TYPE={{ jaeger_agent_span_storage_type }}"
ExecReload=/bin/kill -HUP $MAINPID
ExecStart=/usr/local/bin/jaeger-agent --config-file /etc/jaeger-agent.yml

SyslogIdentifier=jaeger-agent
KillMode=process
Restart=always
RestartSec=5

LockPersonality=true
NoNewPrivileges=true
MemoryDenyWriteExecute=true
PrivateTmp=true
ProtectHome=true
RemoveIPC=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
